# cloudflared-tunnel.service
[Unit]
Description=Cloudflare Tunnel for Dashboard Social
After=network.target
Wants=dashboard-social.service

[Service]
Type=simple
User=root
WorkingDirectory=/data/.openclaw/workspace/dashboard-social/cloudflared-config
ExecStart=/home/linuxbrew/.linuxbrew/bin/cloudflared tunnel --config config.yml run
Restart=always
RestartSec=10
StandardOutput=append:/var/log/cloudflared-tunnel.log
StandardError=append:/var/log/cloudflared-tunnel.log

[Install]
WantedBy=multi-user.target
