# dashboard-social.service
[Unit]
Description=Dashboard Social - Instagram Analytics
After=network.target

[Service]
Type=simple
User=node
WorkingDirectory=/data/.openclaw/workspace/dashboard-social
EnvironmentFile=/data/.openclaw/workspace/dashboard-social/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/dashboard-social.log
StandardError=append:/var/log/dashboard-social.log

[Install]
WantedBy=multi-user.target
